# backup/contrib/init/lightningd_backup-compact.timer
[Unit]
Description=Weekly backup-compact for C-Lightning daemon

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
